# deploy/rp-api.service
[Unit]
Description=Run the ranked-pick-api server

[Service]
ExecStart=/bin/bash /usr/local/bin/rp-api

[Install]
WantedBy=multi-user.target
